# scripts/domain-api.service
[Unit]
Description=Domain Search API
Documentation=https://github.com/moustafamahdy/zonefile-search-tantivy
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/zonefile-search
EnvironmentFile=/opt/zonefile-search/.env
ExecStart=/opt/zonefile-search/target/release/domain-api
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/zonefile-search/data

# Resource limits
MemoryMax=4G
CPUQuota=200%
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
